# systemd/user/qt-scale-update.service
[Unit]
Description=Update Qt scale factor from GNOME monitors.xml

[Service]
Type=oneshot
ExecStart=%h/.local/bin/qt-scale-watch
